# net-fs/netatalk/files/atalk_service.init.2
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/netatalk/files/atalk_service.init.2,v 1.1 2011/10/18 07:37:57 flameeyes Exp $

command="/usr/sbin/${SVCNAME}"

depend() {
	need net atalkd
    use logger dns
}
