# net-fs/netatalk/files/netatalk.init
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-fs/netatalk/files/netatalk.init,v 1.1 2012/09/10 14:28:03 jlec Exp $

command="/usr/sbin/${SVCNAME}"

depend() {
	need net
	need avahi-daemon
	use logger dns
	after entropy
}
