# sys-power/powerman/files/powerman
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-power/powerman/files/powerman,v 1.3 2009/03/04 14:12:22 weaver Exp $

depend() {
	need net
}

start() {
	ebegin "Starting powermand"
	start-stop-daemon --start \
		--pidfile /var/run/powerman/powermand.pid \
		--exec /usr/sbin/powermand
	eend $?
}

stop() {
	ebegin "Stopping powermand"
	start-stop-daemon --stop \
		--pidfile /var/run/powerman/powermand.pid \
		--exec /usr/sbin/powermand
	eend $?
}
