# app-admin/chef-server-webui/files/initd/chef-server-webui
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
	after chef-server
}

CHEF_SERVER_WEBUI_BINARY=${CHEF_SERVER_WEBUI_BINARY:-/usr/bin/chef-server-webui}
CHEF_SERVER_WEBUI_PIDFILE=${CHEF_SERVER_WEBUI_PIDFILE:-/var/run/chef/server-webui.%s.pid}
CHEF_SERVER_WEBUI_LOGFILE=${CHEF_SERVER_WEBUI_LOGFILE:-/var/log/chef/server-webui.log}
CONFIGFILE=${CONFIGFILE:-/etc/chef/webui.rb}
CLUSTER_NODES=${CLUSTER_NODES:-1}
CHEF_USER=${CHEF_USER:-chef}
CHEF_GROUP=${CHEF_GROUP:-chef}

start_pre() {
	checkpath -q -d -m 0750 -o ${CHEF_USER}:${CHEF_GROUP} \
		/etc/chef \
		/etc/chef/certificates \
		/var/{lib,log,run}/chef

	checkpath -q -f -m 0600 -o ${CHEF_USER}:${CHEF_GROUP} \
		${CONFIGFILE}
}

start() {
	ebegin "Starting Chef server web interface"
	${CHEF_SERVER_WEBUI_BINARY} -d -P ${CHEF_SERVER_WEBUI_PIDFILE} --name \
		chef-server-webui -C ${CONFIGFILE} -u ${CHEF_USER} -G ${CHEF_GROUP} \
		-c ${CLUSTER_NODES} -L ${CHEF_SERVER_WEBUI_LOGFILE} ${CHEF_SERVER_WEBUI_OPTS} > /dev/null
	eend $?
}

stop() {
	ebegin "Shutting down Chef server web interface"
	${CHEF_SERVER_WEBUI_BINARY} -K all -P ${CHEF_SERVER_WEBUI_PIDFILE} > /dev/null
	eend $?
}
