# net-wireless/bluez/files/rfcomm-init.d-r1
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-wireless/bluez/files/rfcomm-init.d-r1,v 1.1 2014/03/30 08:21:22 pacho Exp $

depend() {
	after coldplug
	need dbus localmount hostname
}

start() {
	ebegin "Starting rfcomm"
	/usr/bin/rfcomm bind all
	eend $?
}

stop() {
	ebegin "Shutting down rfcomm"
	/usr/bin/rfcomm release all
	eend $?
}
