# net-dns/pdnsd/files/pdnsd.rc6.1
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dns/pdnsd/files/pdnsd.rc6.1,v 1.2 2011/12/31 21:22:49 idl0r Exp $

extra_started_commands="reload"

depend() {
	need net
	provide dns
}

checkconfig() {
	if [ ! -e /etc/pdnsd/pdnsd.conf ] ; then
		eerror "You need an /etc/pdnsd/pdnsd.conf file first."
		eerror "There is a sample in /etc/pdnsd/pdnsd.conf.sample"
		return 1
	fi

	if [ ! -d /var/cache/pdnsd ] ; then
		mkdir /var/cache/pdnsd
		chown pdnsd:pdnsd /var/cache/pdnsd
	fi
}

start() {
	checkconfig || return 1
	ebegin "Starting pdnsd"
	start-stop-daemon --start --exec /usr/sbin/pdnsd \
		-- -s -t -d -p /var/run/pdnsd.pid ${PDNSDCONFIG}
	eend $?
}

stop() {
	ebegin "Stopping pdnsd"
	start-stop-daemon --stop --pidfile /var/run/pdnsd.pid
	eend $?
}

reload() {
	ebegin "Reloading pdnsd.conf"
	/usr/sbin/pdnsd-ctl config
	eend $?
}
