# net-dns/pdnsd/files/pdnsd.online.1
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-dns/pdnsd/files/pdnsd.online.1,v 1.2 2011/03/28 14:00:21 polynomial-c Exp $

depend() {
	need pdnsd
}

wait_for_socket() {
	local i=0 pdnsd_socket="/var/cache/pdnsd/pdnsd.status"
	while [ $i -lt 10 ]; do
		[ -S "${pdnsd_socket}" ] && return 0
		i=$(($i+1))
		sleep 1s
	done

    	return 1
}

start() {
	wait_for_socket || return 1
	ebegin "Starting pdnsd-online"
	/usr/sbin/pdnsd-ctl -q server all up
	eend $?
}

stop() {
	ebegin "Stopping pdnsd-online"
	/usr/sbin/pdnsd-ctl -q server all down
	eend $?
}
