# sys-libs/freeipmi/files/ipmidetectd.initd.3
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-libs/freeipmi/files/ipmidetectd.initd.3,v 1.1 2011/10/18 09:01:08 flameeyes Exp $

command="/usr/sbin/${SVCNAME}"
pidfile=/var/run/$SVCNAME.pid
start_stop_daemon_args="--background --make-pidfile"
