# www-servers/mongrel_cluster/files/mongrel_cluster
#!/sbin/runscript
#
# INIT Script for mongrel_cluster
#
# Robert Sharp Jan 2010

depend() {
  need net
}

start() {
	ebegin "Starting Mongrel Clusters"
	if [ ! -d "$CONF_DIR" ]; then
		eerror "No configuration directory"
		return 1
	fi
	mkdir -p $PID_DIR
	chown $USER:$USER $PID_DIR

	mongrel_cluster_ctl start --clean -c $CONF_DIR

	eend $? "Error starting Mongrel"
}

stop() {
	ebegin "Stopping Mongrel Clusters"
	mongrel_cluster_ctl stop -c $CONF_DIR
	eend $?
}

restart() {
	ebegin "Restarting Mongrel Clusters"
	mongrel_cluster_ctl restart -c $CONF_DIR
	eend $?
	}
	
status() {
	ebegin "Checking Status of Mongrel Clusters"
	mongrel_cluster_ctl status -c $CONF_DIR
	eend $?
	}
	
