# app-admin/ulogd/files/ulogd-0.98
#!/sbin/runscript
#
# chkconfig: - 91 35
# description: Starts and stops the ulogd daemon
#
# config:  /etc/ulogd.conf

extra_started_commands="reload"

depend() {
	need net
	use mysql
}

initService() {
# Avoid using root's TMPDIR
unset TMPDIR

# Check that ulogd.conf exists.
[ -f /etc/ulogd.conf ] || exit 0

RETVAL=0
}

start() {
	initService
	ebegin "Starting ulogd"
	start-stop-daemon --start --quiet --exec /usr/sbin/ulogd -- -d >/dev/null 2>&1
	eend $?
}

stop() {
	initService
	ebegin "Stopping ulogd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/ulogd >/dev/null 2>&1
	eend $?
}

reload() {
	initService
	ebegin "Reloading ulogd.conf file"
	killall -HUP ulogd
	eend $?
}
