# dev-java/openjms/files/0.7.6.1/rc2
#!/sbin/runscript

depend() {
	need net
}

check_config() {
	[ -n "${JMSDIR}" ] || return 1 
	if [ -z "${JAVA_HOME}" ]; then
		JAVA_HOME="$(/usr/bin/java-config --jre-home)"
		if [ -z "${JAVA_HOME}" ]; then
			JAVA_HOME="$(/usr/bin/java-config --jdk-home)"
		fi
	fi
	export JAVA_HOME
	[ -n "${JAVA_HOME}" ] || return 1
}

# See how we were called.
start() {
	check_config || return 1
	ebegin "Starting JMS server"
	cd ${JMSDIR}bin
	./startup.sh > /var/log/openjms.log 2>&1 &
	sleep 5
	if ps ax|grep '[o]rg\.exolab\.jms\.server\.JmsServer' > /dev/null; then
		eend 0
	else
		eend 1
	fi
}
stop() {
	check_config || return 1
	ebegin "Stopping JMS server"
	cd ${JMSDIR}/bin
	if ./shutdown.sh 2>&1 |grep '[Ff]ailed' > /dev/null; then
		eend 1
	else
		eend 0
	fi
}
