# net-misc/lldpd/files/lldpd-initd-1
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/lldpd/files/lldpd-initd-1,v 1.1 2012/02/09 23:34:41 chutzpah Exp $

depend() {
	use net
}
start() {
	ebegin "Starting lldpd"
	start-stop-daemon --start --exec /usr/sbin/lldpd \
	    --pidfile /var/run/lldpd.pid \
	    -- ${LLDPD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping lldpd"
	start-stop-daemon --stop --exec /usr/sbin/lldpd \
	    --pidfile /var/run/lldpd.pid
	eend $?
}
