# dev-db/firebird/files/firebird.init.d.2.5
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-db/firebird/files/firebird.init.d.2.5,v 1.2 2011/05/12 18:42:35 jer Exp $

depend() {
	need net
}

start(){
	ebegin "Starting Firebird server"
	start-stop-daemon --start --pidfile $PIDFILE --user $FBUSER --group $FBGROUP --exec $FBGUARD -- $FB_OPTS
	eend $?
}

stop(){
	ebegin "Stopping Firebird server"
	start-stop-daemon --stop --name fbguard
	eend $?
}

restart(){
	svc_stop
	sleep 1
	svc_start
}

