# x11-misc/keytouch/files/keytouch-acpid
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/x11-misc/keytouch/files/keytouch-acpid,v 1.1 2007/03/21 00:07:50 nyhm Exp $

depend() {
	need acpid
}

start() {
	ebegin "Running keytouch-init"
	/usr/bin/keytouch-init
	eend $?

	ebegin "Starting keytouch-acpid"
	start-stop-daemon --start --background --exec /usr/bin/keytouch-acpid
	eend $?
}

stop() {
	ebegin "Stopping keytouch-acpid"
	start-stop-daemon --stop --exec /usr/bin/keytouch-acpid
	eend $?
}
