# x11-misc/keytouch/files/keytouch
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/x11-misc/keytouch/files/keytouch,v 1.2 2007/01/30 19:52:53 nyhm Exp $

start() {
	ebegin "Running keytouch-init"
	/usr/bin/keytouch-init
	eend $?
}
