# sys-fs/iprutils/files/iprupdate
#!/sbin/runscript

depend() {
	need logger iprinit
}

start() {
	ebegin "Checking ipr microcode levels"
	/sbin/iprupdate &> /dev/null
	eend $? "Failed to complete ipr microcode updates"
}
