# app-office/openerp-web/files/openerp-web
#!/sbin/runscript

# Additional options that are passed to the daemon.
DAEMON_OPTS="-c ${CONFIGFILE}"

depend() {
	after openerp-server
}

start() {
	ebegin "Starting openerp-web"
	start-stop-daemon --start --quiet --pidfile "${PIDFILE}" --exec /usr/bin/openerp-web \
		--make-pidfile --background -- ${DAEMON_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping openerp-web"
	start-stop-daemon --stop --quiet --pidfile "${PIDFILE}" --oknodo
	rm -f "${PIDFILE}"
	eend $?
}
