# net-im/pymsn-t/files/pymsn-t-initd
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-im/pymsn-t/files/pymsn-t-initd,v 1.1 2011/09/04 10:26:35 hanno Exp $

depend() {
	need net
	use jabber-server
}

start() {
	ebegin "Starting MSN Jabber Transport"
	start-stop-daemon --start --pidfile /var/run/jabber/pymsn-t.pid -u jabber -g jabber \
		--exec /usr/bin/python INSPATH/pymsn-t.py -- \
		-b -c /etc/jabber/pymsn-t.xml -l /var/log/jabber/pymsn-t.log
	eend $?
}

stop() {
	ebegin "Stopping MSN Jabber Transport"
	start-stop-daemon --stop --quiet --pidfile /var/run/jabber/pymsn-t.pid
	eend $?
}
