# dev-ruby/stompserver/files/stompserver
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-ruby/stompserver/files/stompserver,v 1.1 2009/03/25 17:02:32 caleb Exp $

start() {
	ebegin "Starting Stompserver"
	start-stop-daemon --start --quiet -b -m -p /var/run/stompserver.pid --exec /usr/bin/stompserver
	eend $?
}

stop() {
	ebegin "Shutting down Stompserver"
	start-stop-daemon --stop --quiet -p /var/run/stompserver.pid
	eend $?
}

