# app-admin/flexlm/files/flexlm-init
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-admin/flexlm/files/flexlm-init,v 1.4 2006/09/26 04:10:49 eradicator Exp $

depend() {
	need net
}

start() {
	ebegin "Starting lmgrd"
	su - ${LMUSER} -c "umask 022; /opt/flexlm/bin/lmgrd -c '${LMLICENSE}' -l '${LMLOG}'"
	eend $?
}

stop() {
	ebegin "Stopping lmgrd"
	/opt/flexlm/bin/lmdown -q -all -force -c "${LMLICENSE}" >& /dev/null
	eend $?
}
