# app-admin/osiris/files/osirisd-4.2.3
#!/sbin/runscript

depend() {
	use logger dns
	need net
}

checkconfig() {
	if [ ! -x /usr/sbin/osirisd ] ; then
		eerror "Daemon image not found"
		return 1
	fi
	return 0
}

start() {
	ebegin "Starting osirisd"
	checkconfig || ( eend $? "Config problems for osirisd daemon" || return 1 )
	start-stop-daemon --start --exec /usr/sbin/osirisd -- -e /var/run/osiris/osirisd.pid ${OSIRISD_OPTS}
	eend $? "Failed to start osirisd daemon"
}

stop() {
	ebegin "Stop osirisd"
	start-stop-daemon --stop --name osirisd --pidfile /var/run/osiris/osirisd.pid
	eend $? "Failed to stop osirisd daemon"
}

