# app-admin/osiris/files/osirismd-4.2.3
#!/sbin/runscript

depend() {
	use logger dns
	need net
}

checkconfig() {
	if [ ! -x /usr/sbin/osirisd ] ; then
		eerror "Daemon image not found"
		return 1
	fi
	return 0
}

start() {
	ebegin "Starting osirismd"
	checkconfig || ( eend $? "Config problems for osirismd daemon" || return 1 )
	start-stop-daemon --start --exec /usr/sbin/osirismd -- -e /var/run/osiris/osirismd.pid ${OSIRISMD_OPTS}
	eend $? "Failed to start osirismd daemon"
}

stop() {
	ebegin "Stop osirismd"
	start-stop-daemon --stop --quiet --pidfile /var/run/osiris/osirismd.pid
	eend $? "Failed to stop osirismd daemon"
}
