# net-analyzer/openvas-scanner/files/openvassd.service
[Unit]
Description=OpenVAS Scanner
After=network.target

[Service]
ExecStart=/usr/sbin/openvassd -f
