# net-misc/teamviewer/files/teamviewer.service
[Unit]
Description=Daemon for TeamViewer-8

[Service]
ExecStart=/opt/teamviewer-8/teamviewerd -f

[Install]
WantedBy=multi-user.target

