# media-plugins/vdr-vdrrip/files/vdrrip-qh
#!/sbin/runscript
# Distributed under the terms of the GNU General Public License, v2 or later
# $Header: /var/cvsroot/gentoo-x86/media-plugins/vdr-vdrrip/files/vdrrip-qh,v 1.2 2006/10/13 21:36:31 hd_brummy Exp $

# default dir's / files
QUEUE="/etc/vdr/plugins/queue.vdrrip"

depend() {
	need localmount
}

start() {
	ebegin "Starting vdrrip queuehandler"
	start-stop-daemon --start --background \
	--make-pidfile \
	--pidfile /var/run/vdrrip-qh.pid --user ${VDRRIP_USER} \
	--exec /usr/bin/vdrrip-qh -- ${QUEUE} ${VDRRIP_TEMPDIR}
	eend $? "Couldn't start vdrrip-qh"
}

stop() {
	ebegin "Stopping vdrrip queuehandler"
	start-stop-daemon --stop --pidfile /var/run/vdrrip-qh.pid \
	--user ${VDRRIP_USER} /usr/bin/vdrrip-qh
	eend $? "Couldn't stop vdrrip-qh."
}
