# app-admin/ulog-acctd/files/init.d/ulog-acctd
#!/sbin/runscript
#
# chkconfig: - 91 35
# description: Starts and stops the ulog-acctd daemon

# We should add "pause" and "resume" actions - Radek Podgorny
#opts="${opts} pause resume"


depend() {
	need net
}


start() {
	ebegin "Starting ulog-acctd"
	start-stop-daemon --start --quiet --exec /usr/sbin/ulog-acctd >/dev/null 2>&1
	eend $?
}	

stop() {
	ebegin "Stopping ulog-acctdd"
	start-stop-daemon --stop --quiet --exec /usr/sbin/ulog-acctd >/dev/null 2>&1
	eend $?
}	
