# app-emulation/cloud-init/files/cloud-final.init
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-emulation/cloud-init/files/cloud-final.init,v 1.1 2014/03/16 21:05:12 prometheanfire Exp $

depend() {
    after cloud-config
    provide cloud-final
}

start() {
    cloud-init modules --mode final
    eend 0
}
