# net-misc/nxserver-freeedition/files/nxserver-2.1.0-init
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# opts="${opts} reload configtest"

NX_SERVER_CMD="/usr/NX/bin/nxserver"

depend() 
{
	need net
	use nxsensor
}

start() {
	ebegin "Starting nxserver"
	$NX_SERVER_CMD --start
	if [ $? = 0 ]; then
		$NX_SERVER_CMD --statistics start
		eend $?
	else
		eend 1
	fi
}

stop() {
	ebegin "Stopping nxserver"
	$NX_SERVER_CMD --stop
	$NX_SERVER_CMD --statistics stop
	eend 0
}
