# dev-db/couchdb/files/couchdb.init-3
#!/sbin/runscript
# Copyright 1999-2013 Dirkjan Ochtman
# Distributed under the terms of the Apache License, Version 2.0

pidfile=${COUCHDB_PID_FILE}
command=${EXEC:-/usr/bin/couchdb}
command_args="-b -o ${COUCHDB_STDOUT_FILE} -e ${COUCHDB_STDERR_FILE} -p ${pidfile} ${COUCHDB_OPTIONS}"

user=${COUCHDB_USER}
group=${COUCHDB_GROUP}

depend() {
	need net
}

start_pre() {
	checkpath -q -d -m 0755 -o couchdb:couchdb /var/run/couchdb
}
