# net-misc/sobby/files/sobby-init-0.4.3
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2         
# $Header: /var/cvsroot/gentoo-x86/net-misc/sobby/files/sobby-init-0.4.3,v 1.1 2007/08/18 23:00:42 dev-zero Exp $

depend() {
	use avahi
}

start() {
	ebegin "Starting obby dedicated server"
	start-stop-daemon -b --quiet --start --exec /usr/bin/sobby -- ${SOBBY_OPTS}
	eend $? "Failed to start sobby"
}

stop() {
	ebegin "Stopping obby dedicated server"
	start-stop-daemon --stop --quiet --exec /usr/bin/sobby
	eend $? "Failed to stop sobby"
}
