# sys-power/cpufrequtils/files/cpufrequtils-init.d-007
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-power/cpufrequtils/files/cpufrequtils-init.d-007,v 1.1 2013/04/21 09:21:06 lxnay Exp $

change() {
	ebegin "Running cpufreq-set ${opts}"
	/usr/libexec/cpufrequtils-change.sh "${@}"
	eend $?
}

start() {
	change "${START_OPTS}" ${SYSFS_EXTRA}
}

stop() {
	change "${STOP_OPTS}"
}
