# app-admin/chef-solr/files/initd/chef-solr
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

CHEF_SOLR_LOGFILE=${CHEF_SOLR_LOGFILE:-/var/log/chef/solr.log}
CONFIGFILE=${CONFIGFILE:-/etc/chef/solr.rb}

pidfile=${CHEF_SOLR_PIDFILE:-/var/run/chef/solr.pid}
command=${CHEF_SOLR_BINARY:-/usr/bin/chef-solr}
command_args="-d -P ${pidfile} -L ${CHEF_SOLR_LOGFILE} -c ${CONFIGFILE} ${CHEF_SOLR_OPTS}"

user=${CHEF_USER:-chef}
group=${CHEF_GROUP:-chef}

depend() {
	need net rabbitmq
	after couchdb
}

start_pre() {
	checkpath -q -d -m 0750 -o ${user}:${group} /etc/chef /var/{lib,log,run}/chef
	checkpath -q -f -m 0600 ${CONFIGFILE}
}
